# crontab
# Run every day at 7pm GMT/Noon PDT
0 19 * * * java -jar /opt/form-flow-starter-app/app.jar checkForUnsubmittedDocs

# Run every day at 5pm PDT (Pausing until production issues are resolved)
# 0 0 * * * java -jar /opt/form-flow-starter-app/app.jar transmit
