# crontab
# Run every day at Noon PDT
0 12 * * * java -jar [build.jar] -Dspring.main.web-application-type=NONE checkForUnsubmittedDocs
